# etc/cron.d/keepalived-master
* * * * * root    [ -f /etc/keepalived/keepalived.conf ] && /var/openpanel/tools/openpanel-keepalived-status > /var/openpanel/http/modules/keepalived/keepalived-status.html
