# etc/cron.d/keepalived-master
* * * * * root    [ -f /etc/keepalived/keepalived.conf ] && /var/openpanel/tools/openpanel-keepalived-status > /var/openpanel/modules/Keepalived.module/html/keepalived-status.html
