# etc/cron.d/keepalived-slave
*/2 * * * * root    [ -f /etc/openpanel/ha/active_slave ] && /usr/sbin/fetch_master_settings
